# wm_i3_sway/src/dot-config/autostart/dunst.desktop
[Desktop Entry]
Exec=/usr/bin/dunst
TryExec=/usr/bin/dunst
Name=Dunst Desktop notifications
Comment=Desktop notifications daemon
Type=Application
