# wm_i3_sway/src/dot-config/autostart/nitrogen.desktop
[Desktop Entry]
Exec=/usr/bin/nitrogen
Name=nitrogen
Comment=X11 background previewer and setter
Type=Application
OnlyShowIn=i3
