# wm_i3_sway/src/dot-config/autostart/org.gnome.Evolution-alarm-notify.desktop
[Desktop Entry]
Exec=/usr/libexec/evolution-data-server/evolution-alarm-notify
Name=evolution-alarm-notify
Type=Application
OnlyShowIn=Gnome
