# gnome/src/dot-config/autostart/ydotool.desktop
[Desktop Entry]
Exec=/usr/bin/ydotoold
TryExec=/usr/bin/ydotoold
Name=ydootol
Comment=Daemon for ydotool
Type=Application
